# deploy/wxbt-bot.service
[Unit]
Description=WXBT telegram bot (long-poll)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO__
ExecStart=__REPO__/.venv/bin/python -u scripts/telegram_bot.py --poll
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
