# hosting/deploy/wxbt-bot.service
[Unit]
Description=WXBT telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/wxbt-hosting
ExecStart=/usr/bin/env bash %h/wxbt-hosting/run_bot.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
